# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

; GPU-Accelerated X Server (MUST be root)
[program:Xorg]
command=/usr/bin/Xorg :0 -noreset -keeptty -verbose 3
priority=10
autostart=true
autorestart=true
user=root
environment=DISPLAY=":0"
stdout_logfile=/tmp/xorg.log
stderr_logfile=/tmp/xorg_err.log

; Audio Server
[program:pulseaudio]
command=/usr/bin/pulseaudio --disallow-exit --exit-idle-time=-1 --log-target=stderr
priority=15
autostart=true
autorestart=true
user=gamer
environment=DISPLAY=":0",HOME="/home/gamer"
stdout_logfile=/tmp/pulse.log
stderr_logfile=/tmp/pulse_err.log

; Sunshine Game Streaming Server
[program:sunshine]
command=/usr/bin/sunshine
priority=20
autostart=true
autorestart=true
user=gamer
environment=DISPLAY=":0",HOME="/home/gamer",XAUTHORITY="/home/gamer/.Xauthority",PULSE_SERVER="unix:/run/user/1000/pulse/native"
stdout_logfile=/tmp/sunshine.log
stderr_logfile=/tmp/sunshine_err.log
startsecs=3
